# docker/supervisor/scheduler.conf
[program:scheduler]
command=php /var/www/html/bin/console scheduler:consume
numprocs=1
autostart=true
autorestart=true
user=www-data